# templates/default/hiveserver2.service.erb
[Unit]
Description = Hive metaserver
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=120
<% if node['hops']['alive-watchdog']['enabled'].casecmp?("true") -%>
# Do not be very strict as it is normal for the service to fail and restart immediatelly
# when liveness watchdog is enabled
StartLimitBurst=350
<% else -%>
StartLimitBurst=18
<% end -%>


[Service]
User = <%= node['hive2']['user'] %>

ExecStart = <%= node['hive2']['base_dir'] %>/bin/hiveserver2
LimitNOFILE = 20000
TimeoutSec= 10min
Restart=on-failure
RestartSec=2

Type = simple

[Install]
WantedBy = multi-user.target
